# local/etc/periodic/weekly/120.cabal-maintenence
#!/bin/sh
# Perform some regular maintinence of Cabal

rc=0
/Users/wren/.cabal/bin/cabal update || rc=3
exit $rc
