# local/etc/periodic/monthly/200.gpp-refresh-keys
#!/bin/sh
echo ''
echo 'Updating GnuPG key database:'

rc=0
gpg --refresh-keys || rc=3
exit $rc
